# systemd/piboard-kiosk.service
[Unit]
Description=PiBoard Chromium Kiosk
After=network-online.target systemd-user-sessions.service
Wants=network-online.target

[Service]
Type=simple
User=piboard
Group=piboard
WorkingDirectory=/opt/piboard-kiosk
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/opt/piboard-kiosk/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/piboard-kiosk/scripts/start-kiosk-x.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
